# apps/zsp-aitool/deploy/systemd/zsp-hyperframes-worker.service
[Unit]
Description=zsp-aitool HyperFrames Render Worker (optional)
After=network.target

[Service]
User=zeazdev
WorkingDirectory=/home/zeazdev/zsp-aitool
EnvironmentFile=/home/zeazdev/zsp-aitool/.env
ExecStart=/usr/bin/npm run hyperframes:worker
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/zsp-aitool/hyperframes

[Install]
WantedBy=multi-user.target
