# apps/zsp-aitool/deploy/systemd/zsp-hyperframes-cleanup.service
[Unit]
Description=zsp-aitool HyperFrames Render Artifact Cleanup
After=network.target

[Service]
Type=oneshot
User=zeazdev
WorkingDirectory=/home/zeazdev/zsp-aitool
EnvironmentFile=/home/zeazdev/zsp-aitool/.env
ExecStart=/usr/bin/npm run hyperframes:cleanup-renders
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/zsp-aitool/hyperframes
